# roles/marathon/templates/marathon.conf.j2
description "Marathon scheduler for Mesos"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

export LIBPROCESS_IP={{ marathon_libprocess_ip }}

exec /usr/bin/marathon
